# supervisor-example.conf
[program:tank-temp]
command=/home/robin/temp/run.sh
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/home/robin/temp/temp.log
directory=/home/robin/temp
